# deployment/gatorweb.conf
description "The website Delegators use to fulfill orders"

start on runlevel [2345]
stop on runlevel [!2345]

setuid ubuntu
setgid www-data

chdir /var/gator/web/
exec /var/gator/web/startserver.py
